# systemd.service
[Unit]
Description=DAEMON_NAME

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=DAEMON_PATH

[Install]
WantedBy=multi-user.target
